# extensions/job-hunt/launchd/com.openbrain.posting-maintenance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.openbrain.posting-maintenance</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/dfrysinger/Projects/open-brain/extensions/job-hunt/scripts/launchd-wrapper.sh</string>
    <string>/Users/dfrysinger/Projects/open-brain/extensions/job-hunt/scripts/posting-maintenance.ts</string>
    <string>--mode</string>
    <string>check-active</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>0</integer>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/posting-maintenance.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/posting-maintenance.err</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/dfrysinger</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
